# tools/com.johnherrick.cascadia.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.johnherrick.cascadia.dashboard</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/johnherrick/cascadia/target/release/cascadia-api</string>
    <string>--listen</string>
    <string>0.0.0.0:5187</string>
    <string>--static-dir</string>
    <string>apps/web/dist</string>
    <string>--cluster-history-path</string>
    <string>artifacts/cluster/telemetry-v1.jsonl</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/johnherrick/cascadia</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>10</integer>
  <key>ProcessType</key>
  <string>Interactive</string>
  <key>StandardOutPath</key>
  <string>/Users/johnherrick/cascadia/artifacts/logs/dashboard.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/johnherrick/cascadia/artifacts/logs/dashboard.err.log</string>
</dict>
</plist>
