# tools/com.johnherrick.cascadia.adr0078-supervisor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.johnherrick.cascadia.adr0078-supervisor</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/johnherrick/cascadia/.venv/bin/python</string>
    <string>/Users/johnherrick/cascadia/tools/adr0078_cluster_supervisor.py</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/johnherrick/cascadia</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <false/>
  <key>ProcessType</key>
  <string>Background</string>
  <key>StandardOutPath</key>
  <string>/Users/johnherrick/cascadia/artifacts/logs/adr0078-launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/johnherrick/cascadia/artifacts/logs/adr0078-launchd.err.log</string>
</dict>
</plist>
